# src/terok_dbus/resources/systemd/terok-dbus.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# systemd user service for the terok-dbus Shield1 clearance hub.
#
# The {{BIN}} placeholder is resolved at install time by the install
# command (e.g. via sys.executable -m terok_dbus._cli) so the unit
# works regardless of PATH — matching the pattern used by
# terok-credential-proxy and terok-gate.
#
# Install:
#   terok-dbus install-service
#   — or manually —
#   sed "s|{{BIN}}|$(which terok-dbus)|" contrib/systemd/terok-dbus.service \
#       > ~/.config/systemd/user/terok-dbus.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-dbus
#
# Logs:
#   journalctl --user -u terok-dbus -f

[Unit]
Description=terok Shield1 clearance hub (verdict dispatcher + desktop notifier)
Documentation=https://github.com/terok-ai/terok-dbus
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}} serve
Restart=on-failure
RestartSec=5

# ── Hardening: filesystem guards only ─────────────────────────────────
#
# The hub's threat model is "process-level compromise via a dbus-fast
# parsing bug or malformed JSON from a container reader".  Constraining
# the hub's filesystem view shrinks what such a compromise can reach —
# an attacker sitting in the hub's memory can't read SSH keys, GPG
# state, or shell history if $HOME is read-only; can't tamper with
# installed binaries if /usr is read-only.  Worth having.
#
# None of the directives in this block trigger the NNP / SELinux
# cascade that bit us with the seccomp-based family
# (`RestrictAddressFamilies`, `SystemCallFilter`, etc. — see the
# *rejected* block below).
ProtectHome=read-only
# ProtectHome=read-only covers /run/user/$UID too, which is where the
# event ingester needs to create and unlink its unix socket.  %t
# expands to $XDG_RUNTIME_DIR for user units — this carves exactly the
# path the ingester owns back into writable.
ReadWritePaths=%t
ProtectSystem=full

# ── Rejected: anything that compiles to a seccomp filter ──────────────
#
# The kernel refuses to install a seccomp filter on an unprivileged
# process unless `no_new_privs` is set (see `seccomp(2)`).  Once NNP
# is on, SELinux refuses the `unconfined_t → container_runtime_t`
# transition that `podman unshare` does on every verdict apply — the
# user sees "Allow failed: X" on every click with an `nnp_transition`
# AVC in the audit log.  `NoNewPrivileges=no` in the unit does NOT win
# that fight: systemd has to enable NNP to make the seccomp filter
# load at all.
#
# Directives that compile to seccomp and therefore cannot be used:
#
#   RestrictAddressFamilies=
#   SystemCallFilter=  SystemCallArchitectures=  SystemCallLog=
#   MemoryDenyWriteExecute=yes
#   LockPersonality=yes
#   RestrictNamespaces=
#   RestrictRealtime=yes
#   RestrictSUIDSGID=yes
#   NoExecPaths=  ExecPaths=
#
# Directives that *explicitly* imply NNP (per systemd.exec(5)) and
# therefore also cannot be used:
#
#   PrivateDevices=yes
#   ProtectKernelTunables=yes  ProtectKernelModules=yes
#   ProtectKernelLogs=yes  ProtectClock=yes
#   DynamicUser=yes  ProtectHostname=yes
#
# A separate ``terok-verdict`` helper service (that does the
# ``terok-shield allow|deny`` exec without NNP, talking to this hub
# over a local socket) would recover this hardening — see the
# deferred plan in the terok-dbus design notes.  That's the proper
# long-term fix; for now the filesystem guards above are what we can
# do without splitting the process tree.

[Install]
WantedBy=default.target
